# websockets.conf
[program:websockets]
command=/usr/bin/php /www/wwwroot/gariin.test/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true